# net-analyzer/splunk-bin/files/splunk.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

start() {
	ebegin "Starting Splunk"
	/opt/splunk/bin/splunk start
	eend $?
}

stop() {
	ebegin "Stopping Splunk"
	/opt/splunk/bin/splunk stop
	eend $?
}
